# systemd/docker-watch-heartbeat.timer
[Unit]
Description=Run Docker Watch heartbeat daily

[Timer]
OnCalendar=*-*-* 09:05:00
Persistent=true

[Install]
WantedBy=timers.target
