# systemd/docker-watch.service
[Unit]
Description=Docker Watch (events -> Telegram)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=dockwatch
Group=dockwatch
SupplementaryGroups=docker

ExecStart=/opt/docker-watch/docker-watch.sh

RuntimeDirectory=docker-watch
RuntimeDirectoryMode=0750

# Hardening baseline
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true

RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# En Ubuntu /var/run suele ser symlink de /run
ReadWritePaths=/run/docker-watch /var/run/docker.sock /run/docker.sock

Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
